# packaging/systemd/ecallisto-acquire.service.in
[Unit]
Description=e-Callisto NG acquisition daemon (scheduler + uploader)
After=network-online.target chrony.service
Wants=network-online.target

[Service]
Type=simple
User=@USER@
WorkingDirectory=@APP_DIR@
EnvironmentFile=@ENV_FILE@
# Acquisition owns the scheduler/uploader loops; the web service runs with
# run_loops_in_web=false so they are not double-run.
ExecStart=@VENV@/bin/ecallisto-ng acquire
Restart=always
RestartSec=3
# Never lose un-uploaded science data on restart.
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
